# deploy/zztmmo.service
[Unit]
Description=ZZTMMO browser test server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/zztmmo
EnvironmentFile=/opt/zztmmo/.env
ExecStart=/opt/zztmmo/zzt-server -addr 127.0.0.1:8080 -world TOWN -web web/dist -help . -saves saves
Restart=on-failure
RestartSec=2
# On stop/restart, zzt-server catches SIGTERM and warns connected players for
# -shutdown-grace (60s) so they can save before it exits. TimeoutStopSec must
# exceed that grace + the ~5s HTTP drain, or systemd would SIGKILL mid-warning.
TimeoutStopSec=80

[Install]
WantedBy=multi-user.target
